# EdiabasLib/EnetAdapter/files/etc/uci-defaults/01_set_passwd
#!/bin/sh
lua /usr/bin/set_passwd root root
exit 0
