# EdiabasLib/EnetAdapter/files/etc/uci-defaults/03_network
#!/bin/sh
uci set network.lan.ipaddr='192.168.100.1'
uci commit network
exit 0
