# EdiabasLib/EnetAdapter/files/etc/uci-defaults/luci-theme-z-default
uci set luci.main.mediaurlbase=/luci-static/bootstrap
uci commit luci

exit 0
