# EdiabasLib/EnetAdapter/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

swconfig dev switch0 port 1 set disable 1
swconfig dev switch0 port 2 set disable 1
swconfig dev switch0 port 3 set disable 1
swconfig dev switch0 port 4 set disable 1
swconfig dev switch0 set apply

exit 0
